# files/usr/lib/systemd/system/mnt-data.mount
[Unit]
Description=cni Directory

[Mount]
What=/dev/disk/by-label/data
Where=/mnt/data
Type=auto
Options=defaults,relatime
