# files/usr/lib/systemd/system/containerd.service.d/10-mod.conf
[Unit]
After=etc-cni.mount var-lib-opt-cni.mount var-lib-containerd.mount
Requires=etc-cni.mount var-lib-containerd.mount

[Service]
ExecStartPre=/bin/bash -exc "mkdir -p /var/run/containerd-shims/ && source /etc/os-release && rsync --ignore-existing -rav /usr/local/share/containerd-shims/ /var/run/containerd-shims/$BUILD_ID/; ln -sf /var/run/containerd-shims/$BUILD_ID /var/run/containerd-shims/latest"
