# files/usr/lib/systemd/system/zfs-zed.service.d/10-config.conf
[Service]
ExecStart=
ExecStart=/usr/bin/zed -F -p /var/run/zed.pid -s /var/run/zed.state
